# newyorktimes-elt/cron/crontab
# Run every day at 12 AM
*/2 * * * * root /app/scheduler.sh >> /var/log/cron.log 2>&1
