# newyorktimes-elt/cron/crontab
# Run every day at 12 AM
0 0 * * * root /app/scheduler.sh >> /var/log/cron.log 2>&1
